# srcs/influxdb/supervisord.conf
  
[supervisord]
nodemon=true





[program:influxd]
command=influxd -config /etc/influxdb.conf
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true


[program:telegraf]
command=telegraf
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true